# .backup-rclone/patients-vacuum.service
[Unit]
Description=Run patient's vacuum once

[Service]
Type=oneshot
ExecStart=/usr/bin/bash %h/dev/dermatodoc/patients/.backup-rclone/vacuum.sh
