# .backup/patients-backup.service
[Unit]
Description=Run patient's backup once

[Service]
Type=oneshot
ExecStart=/usr/bin/bash %h/.backup/backup.sh
