# .backup/patients-backup.service
[Unit]
Description=Run patient's backup once

[Service]
Type=oneshot
ExecStart=/usr/bin/sh %h/dev/dermatodoc/patients/.backup/backup.sh
