# .backup/patients-vacuum.service
[Unit]
Description=Run patient's vacuum once

[Service]
Type=oneshot
ExecStart=/usr/bin/bash %h/.backup/vacuum.sh
